# src/main/resources/external/event-api.service
[Unit]
Description=event-api: Event API

[Service]
Type=forking

# WARNING: environment.rc and config files can override any variable defined here.
# WARNING: Do not guess, but read wrapper.sh script instead.

Environment=APP_JAVA_OPTS="-server -Xms4g -Xmx6g -XX:+HeapDumpOnOutOfMemoryError -Xlog:gc*:file=%h/log/event-api.GC.log"
Environment=CLASSPATH=""
Environment=APP_ENTRY_POINT="-jar %JAR_FILE%"
Environment=APP_ARGS=""
Environment=SPRING_CONFIG_ADDITIONAL_LOCATION="file:%h/config.local.yaml"
Environment=APP_HOME="%h/current"
Environment=APP_BASE="%h"

WorkingDirectory=%h/current
ExecStart=/bin/bash -c "\"%h/bin/event-api\" start"
ExecStop=/bin/bash  -c "\"%h/bin/event-api\" stop"

[Install]
WantedBy=default.target
